# nodes/orchestration/package/service/x86/ubuntu14/nano_agent.conf
description     "CP_Nano_Linux_Agent"
author          "NGEN team, Check Point Technologies LTD."

start on filesystem

respawn

exec /etc/cp/watchdog/cp-nano-watchdog
